# services/signalize-mplayer.service
[Unit]
Description=Signalize MPlayer Service
After=network.target
StartLimitIntervalSec=0
[Service]
Type=simple
Restart=always
RestartSec=1
User=signalize
WorkingDirectory=/home/signalize/core
ExecStart=/usr/bin/php artisan socket mplayer

[Install]
WantedBy=multi-user.target